# relay/systemd/enclave-tcp-relay.service
[Unit]
Description=Enclave TCP relay (SNI router; TLS terminates in-enclave)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/nan-relay/tcp-relay.env
WorkingDirectory=/opt/nan-relay
ExecStart=/usr/bin/node /opt/nan-relay/relay.js
Restart=always
RestartSec=2
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# RELAY_PORTS=1-19999 binds ~20k listen sockets + one fd per live connection
LimitNOFILE=65536
# lets the range include ports below 1024
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
