# relay/systemd/nan-api-relay.service
[Unit]
Description=Enclave API relay (fleet discovery + placement)
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/nan-relay/api-relay.env
WorkingDirectory=/opt/nan-relay
ExecStart=/usr/bin/node /opt/nan-relay/api-relay.js
Restart=always
RestartSec=2
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
